# code/scripts/setup/drai_camera.service
[Unit]
Description=The Camera Webserver of the